# crontab
# Cleanup any completed torrents that have been post processed and are just left seeding...
20 * * * * /home/luca/scripts/remove-completed-downloads.sh
